# etc/crontab
0,15,30,45 * * * * /home/flooddata/scripts/cron.sh
